# 99-auto-mount-drives.rules
KERNEL=="sd[a-z][0-9]", ACTION=="add", RUN+="/usr/bin/systemctl start auto-mount-drives.service"